# DEBIAN/postinst
#!/bin/bash
set -e

# Enable and start the timer
systemctl daemon-reload
systemctl enable cf-nginx.timer
systemctl start cf-nginx.timer

# Create cache directory
mkdir -p /var/lib/cf-nginx

echo "✅ Cloudflare NGINX IPs package installed"
echo "   Run: cf-nginx to configure sites"
echo "   Auto-update: enabled (daily at 3am)"

exit 0
