# DEBIAN/postinst
#!/bin/bash
mkdir -p /var/lib/cf-nginx
touch /var/lib/cf-nginx/enabled_sites
systemctl daemon-reload
systemctl enable cf-nginx.timer
systemctl start cf-nginx.timer
/bin/bash -c 'source /usr/lib/cf-nginx/cf-ips.sh && fetch_cf_ips' 2>/dev/null || true
echo "✅ cf-nginx installed - Run 'cf-nginx' to get started"
exit 0
