# DEBIAN/prerm
#!/bin/bash
systemctl stop cf-nginx.timer 2>/dev/null || true
systemctl disable cf-nginx.timer 2>/dev/null || true
exit 0
